# setup/dataset/services/kiwipanel.service
[Unit]
Description=Kiwipanel application
[Service]
Type=simple
Restart=always
RestartSec=5s
WorkingDirectory=/home/scaffolding/bin
ExecStart=/home/scaffolding/bin/kiwipanel start
[Install]
WantedBy=multi-user.target